# qr3_62.conf
[program:qr3_62]
command=python3 /home/curator3/qr3_62/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator3/qr3_62/test.err.log
stdout_logfile=/home/curator3/qr3_62/test.err.log
